# systemd/user/forge-heartbeat.timer
[Unit]
Description=ForgeOS heartbeat — every 5 minutes
Documentation=https://github.com/VibeCodingLabs/forge-os

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
