# systemd/user/forge-observer.service
[Unit]
Description=ForgeOS observability snapshot

[Service]
Type=oneshot
ExecStart=%h/.local/bin/forge-observer
